# .config/autostart/xfce4-clipman-plugin-autostart.desktop
[Desktop Entry]
Hidden=true
TryExec=xfce4-clipman
Exec=xfce4-clipman

